# config/shell/kubectl-proxy.conf
description "kubectl proxy"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec sudo -u workshop /usr/local/bin/kubectl proxy --accept-hosts="^.*$"